# Robot/Server/cron.hourly
#!/bin/bash
SVNROOT=/svn/svnroot
chgrp acc /data/ -R
chmod 777 /data/ -R
pushd .
cd /data/regnskaber
git add .
git commit -m "autocommit"
git push
source /svn/svnroot/aliases
tpath=~/regnskaber/transactions_stuff
cd $tpath
l b
popd
((ls ~/regnskaber/*/.tags/*|grep -v .diff)|while read line; do basename "$line"; done)|uniq -u > /home/joo/regnskaber/stuff/.tagz.list
